# src/ntpmon-telegraf.upstart
# ntpmon-telegraf - gather NTP statistics and send to telegraf TCP listener

description "NTPmon telegraf collector"

start on started telegraf
start on started ntp

stop on runlevel [06]

respawn

exec {{ install_dir }}/ntpmon.py --interval 60 --mode telegraf
